# tcz/bluez.tcz.list
/usr/local/etc/udev/rules.d/97-bluetooth.rules
/usr/local/etc/alsa/bluetooth.conf
/usr/local/etc/dbus-1/system.d/bluetooth.conf
/usr/local/include/bluetooth/hidp.h
/usr/local/include/bluetooth/cmtp.h
/usr/local/include/bluetooth/bnep.h
/usr/local/include/bluetooth/rfcomm.h
/usr/local/include/bluetooth/sdp_lib.h
/usr/local/include/bluetooth/sdp.h
/usr/local/include/bluetooth/l2cap.h
/usr/local/include/bluetooth/sco.h
/usr/local/include/bluetooth/hci_lib.h
/usr/local/include/bluetooth/hci.h
/usr/local/include/bluetooth/bluetooth.h
/usr/local/sbin/hciemu
/usr/local/sbin/hciconfig
/usr/local/sbin/hciattach
/usr/local/sbin/bluetoothd
/usr/local/bin/pand
/usr/local/bin/rctest
/usr/local/bin/l2test
/usr/local/bin/ciptool
/usr/local/bin/sdptool
/usr/local/bin/hcitool
/usr/local/bin/l2ping
/usr/local/bin/rfcomm
/usr/local/bin/agent
/usr/local/bin/simple-agent
/usr/local/bin/simple-service
/usr/local/lib/libbluetooth.la
/usr/local/lib/libbluetooth.so
/usr/local/lib/libbluetooth.so.3
/usr/local/lib/libbluetooth.so.3.6.0
/usr/local/lib/gstreamer-0.10/libgstbluetooth.la
/usr/local/lib/gstreamer-0.10/libgstbluetooth.so
/usr/local/lib/cups/backend/bluetooth
/usr/local/lib/alsa-lib/libasound_module_ctl_bluetooth.la
/usr/local/lib/alsa-lib/libasound_module_ctl_bluetooth.so
/usr/local/lib/alsa-lib/libasound_module_pcm_bluetooth.la
/usr/local/lib/alsa-lib/libasound_module_pcm_bluetooth.so
/usr/local/lib/pkgconfig/bluez.pc
/usr/local/share/bluez/files/rfcomm.conf
/usr/local/share/bluez/files/main.conf
/usr/local/tce.installed/bluez
/usr/local/etc/init.d/bluez
