# tcz/OSS.tcz.list
./usr/lib/oss/soundon.user
./usr/lib/oss/lib/libflashsupport_32-ssl.so
./usr/lib/oss/lib/libOSSlib.so
./usr/lib/oss/lib/flashsupport.c
./usr/lib/oss/lib/libflashsupport_32.so
./usr/lib/oss/include/sys/ossddk/ossddk.h
./usr/lib/oss/include/sys/ossddk/oss_exports.h
./usr/lib/oss/include/sys/soundcard.h
./usr/lib/oss/include/sys/oss_userdev_exports.h
./usr/lib/oss/include/internals/midi_core.h
./usr/lib/oss/include/internals/timestamp.h
./usr/lib/oss/include/internals/audio_core.h
./usr/lib/oss/include/internals/WARNING.txt
./usr/lib/oss/include/internals/mixer_core.h
./usr/lib/oss/include/internals/local_config.h
./usr/lib/oss/include/midiparser.h
./usr/lib/oss/conf.sample/oss_audigyls.conf
./usr/lib/oss/conf.sample/oss_usb.conf
./usr/lib/oss/conf.sample/oss_als4k.conf
./usr/lib/oss/conf.sample/oss_ymf7xx.conf
./usr/lib/oss/conf.sample/oss_fmedia.conf
./usr/lib/oss/conf.sample/oss_trident.conf
./usr/lib/oss/conf.sample/oss_sbxfi.conf
./usr/lib/oss/conf.sample/oss_emu10k1x.conf
./usr/lib/oss/conf.sample/oss_digi32.conf
./usr/lib/oss/conf.sample/oss_midiloop.conf
./usr/lib/oss/conf.sample/oss_envy24ht.conf
./usr/lib/oss/conf.sample/oss_hdaudio.conf
./usr/lib/oss/conf.sample/oss_userdev.conf
./usr/lib/oss/conf.sample/oss_sbpci.conf
./usr/lib/oss/conf.sample/oss_ich.conf
./usr/lib/oss/conf.sample/oss_sblive.conf
./usr/lib/oss/conf.sample/oss_als3xx.conf
./usr/lib/oss/conf.sample/oss_audioloop.conf
./usr/lib/oss/conf.sample/oss_cs461x.conf
./usr/lib/oss/conf.sample/oss_imux.conf
./usr/lib/oss/conf.sample/oss_allegro.conf
./usr/lib/oss/conf.sample/osscore.conf
./usr/lib/oss/conf.sample/oss_envy24.conf
./usr/lib/oss/scripts/showprocs.sh
./usr/lib/oss/scripts/setup-alsa.sh
./usr/lib/oss/scripts/killprocs.sh
./usr/lib/oss/scripts/oss_usb-create-devices
./usr/lib/oss/scripts/remove_drv.sh
./usr/lib/oss/scripts/restore_drv.sh
./usr/lib/oss/cuckoo/cuckoo_pcm.c
./usr/lib/oss/cuckoo/Makefile
./usr/lib/oss/cuckoo/cuckoo_midi.c
./usr/lib/oss/cuckoo/cuckoo.c
./usr/lib/oss/cuckoo/cuckoo_mixer.c
./usr/lib/oss/cuckoo/checksum.h
./usr/lib/oss/cuckoo/cuckoo.h
./usr/lib/oss/version.dat
./usr/lib/oss/etc.sample/S89oss
./usr/lib/oss/etc.sample/userdefs
./usr/lib/oss/etc.sample/devices.list
./usr/bin/ossinfo
./usr/bin/osstest
./usr/bin/ossplay
./usr/bin/ossmix
./usr/bin/osspartysh
./usr/bin/ossrecord
./usr/sbin/ossdetect
./usr/sbin/soundoff
./usr/sbin/ossdevlinks
./usr/sbin/ossvermagic
./usr/sbin/soundon
./usr/sbin/vmixctl
./usr/sbin/savemixer
./usr/local/tce.menu/OSS
./usr/local/tce.installed/OSS
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_usb.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_sblive.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_ich.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_via823x.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_fmedia.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_neomagic.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_ali5455.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_userdev.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_emu10k1x.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_vortex.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_envy24.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_sbpci.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_maestro.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_audioloop.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_cmi878x.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_midimix.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_cs461x.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_ymf7xx.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_via97.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_sbxfi.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_solo.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_atiaudio.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_audigyls.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_digi96.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_hdaudio.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_midiloop.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_digi32.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_cmpci.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_als3xx.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_trident.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_imux.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_envy24ht.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/osscore.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_audiopci.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_cs4281.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_s3vibes.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_geode.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_allegro.ko
./lib/modules/2.6.29.1-tinycore/kernel/oss/oss_als4k.ko
./etc/asound.conf
./etc/init.d/oss
./etc/oss.conf
