# tcz/OSS40.tcz.list
etc/init.d/oss
etc/oss.conf
lib/modules/2.6.29.1-tinycore/kernel/oss/ali5455.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/allegro.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/als300.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/als4000.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/apci97.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/atiaudio.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/audigyls.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/audioloop.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/audiopci.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/cmi8788.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/cmpci.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/cs4280.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/cs4281.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/digi32.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/digi96.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/emu10k1x.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/envy24.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/envy24ht.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/fm801.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/geode.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/hdaudio.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/ich.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/imux.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/maestro.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/neomagic.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/osscore.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/ossusb.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/s3vibes.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/sblive.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/sbxfi.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/softoss.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/solo.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/trident.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/via8233.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/via97.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/vortex.ko
lib/modules/2.6.29.1-tinycore/kernel/oss/ymf7xx.ko
usr/bin/ossinfo
usr/bin/ossmix
usr/bin/ossplay
usr/bin/ossrecord
usr/bin/osstest
usr/bin/ossxmix
usr/lib/oss/build/.params
usr/lib/oss/build/Makefile.osscore
usr/lib/oss/build/Makefile.tmpl
usr/lib/oss/build/Module.symvers
usr/lib/oss/build/ali5455.c
usr/lib/oss/build/allegro.c
usr/lib/oss/build/als300.c
usr/lib/oss/build/als4000.c
usr/lib/oss/build/apci97.c
usr/lib/oss/build/atiaudio.c
usr/lib/oss/build/audigyls.c
usr/lib/oss/build/audioloop.c
usr/lib/oss/build/audiopci.c
usr/lib/oss/build/build.list
usr/lib/oss/build/cmi8788.c
usr/lib/oss/build/cmpci.c
usr/lib/oss/build/cs4280.c
usr/lib/oss/build/cs4281.c
usr/lib/oss/build/digi32.c
usr/lib/oss/build/digi96.c
usr/lib/oss/build/emu10k1x.c
usr/lib/oss/build/envy24.c
usr/lib/oss/build/envy24ht.c
usr/lib/oss/build/fm801.c
usr/lib/oss/build/geode.c
usr/lib/oss/build/hdaudio.c
usr/lib/oss/build/ich.c
usr/lib/oss/build/imux.c
usr/lib/oss/build/install.sh
usr/lib/oss/build/install.sh.orig
usr/lib/oss/build/maestro.c
usr/lib/oss/build/module.inc
usr/lib/oss/build/modules.order
usr/lib/oss/build/neomagic.c
usr/lib/oss/build/osscore.c
usr/lib/oss/build/osscore.c.old
usr/lib/oss/build/osscore.c.orig
usr/lib/oss/build/osscore_symbols.inc
usr/lib/oss/build/ossddk
usr/lib/oss/build/ossdip.h
usr/lib/oss/build/ossusb.c
usr/lib/oss/build/pci_wrapper.inc
usr/lib/oss/build/regparm
usr/lib/oss/build/s3vibes.c
usr/lib/oss/build/sblive.c
usr/lib/oss/build/sbxfi.c
usr/lib/oss/build/softoss.c
usr/lib/oss/build/solo.c
usr/lib/oss/build/soundcard.h
usr/lib/oss/build/trident.c
usr/lib/oss/build/udi.h
usr/lib/oss/build/usb_wrapper.inc
usr/lib/oss/build/via8233.c
usr/lib/oss/build/via97.c
usr/lib/oss/build/vortex.c
usr/lib/oss/build/wrap.h
usr/lib/oss/build/ymf7xx.c
usr/lib/oss/conf/ac97.conf
usr/lib/oss/conf/allegro.conf
usr/lib/oss/conf/als300.conf
usr/lib/oss/conf/als4000.conf
usr/lib/oss/conf/apci97.conf
usr/lib/oss/conf/audigyls.conf
usr/lib/oss/conf/audioloop.conf
usr/lib/oss/conf/cs4280.conf
usr/lib/oss/conf/digi32.conf
usr/lib/oss/conf/emu10k1x.conf
usr/lib/oss/conf/envy24.conf
usr/lib/oss/conf/envy24ht.conf
usr/lib/oss/conf/fm801.conf
usr/lib/oss/conf/geode.conf
usr/lib/oss/conf/hdaudio.conf
usr/lib/oss/conf/ich.conf
usr/lib/oss/conf/imux.conf
usr/lib/oss/conf/ossusb.conf
usr/lib/oss/conf/sblive.conf
usr/lib/oss/conf/softoss.conf
usr/lib/oss/conf/trident.conf
usr/lib/oss/conf/ymf7xx.conf
usr/lib/oss/cuckoo/Makefile
usr/lib/oss/cuckoo/checksum.h
usr/lib/oss/cuckoo/cuckoo.c
usr/lib/oss/cuckoo/cuckoo.h
usr/lib/oss/cuckoo/cuckoo_midi.c
usr/lib/oss/cuckoo/cuckoo_mixer.c
usr/lib/oss/cuckoo/cuckoo_pcm.c
usr/lib/oss/etc/S89oss
usr/lib/oss/etc/devices.list
usr/lib/oss/etc/installed_drivers
usr/lib/oss/etc/legacy_devices
usr/lib/oss/etc/userdefs
usr/lib/oss/include/internals/WARNING.txt
usr/lib/oss/include/internals/audio_core.h
usr/lib/oss/include/internals/midi_core.h
usr/lib/oss/include/internals/mixer_core.h
usr/lib/oss/include/internals/timestamp.h
usr/lib/oss/include/sys/ossddk/oss_exports.h
usr/lib/oss/include/sys/ossddk/ossddk.h
usr/lib/oss/include/sys/soundcard.h
usr/lib/oss/lib/flashsupport.c
usr/lib/oss/lib/libflashsupport-normal.so
usr/lib/oss/lib/libflashsupport-ssl.so
usr/lib/oss/modules
usr/lib/oss/modules.noregparm/ali5455.o
usr/lib/oss/modules.noregparm/allegro.o
usr/lib/oss/modules.noregparm/als300.o
usr/lib/oss/modules.noregparm/als4000.o
usr/lib/oss/modules.noregparm/apci97.o
usr/lib/oss/modules.noregparm/atiaudio.o
usr/lib/oss/modules.noregparm/audigyls.o
usr/lib/oss/modules.noregparm/audioloop.o
usr/lib/oss/modules.noregparm/audiopci.o
usr/lib/oss/modules.noregparm/cmi8788.o
usr/lib/oss/modules.noregparm/cmpci.o
usr/lib/oss/modules.noregparm/cs4280.o
usr/lib/oss/modules.noregparm/cs4281.o
usr/lib/oss/modules.noregparm/digi32.o
usr/lib/oss/modules.noregparm/digi96.o
usr/lib/oss/modules.noregparm/emu10k1x.o
usr/lib/oss/modules.noregparm/envy24.o
usr/lib/oss/modules.noregparm/envy24ht.o
usr/lib/oss/modules.noregparm/fm801.o
usr/lib/oss/modules.noregparm/geode.o
usr/lib/oss/modules.noregparm/hdaudio.o
usr/lib/oss/modules.noregparm/ich.o
usr/lib/oss/modules.noregparm/imux.o
usr/lib/oss/modules.noregparm/maestro.o
usr/lib/oss/modules.noregparm/neomagic.o
usr/lib/oss/modules.noregparm/ossusb.o
usr/lib/oss/modules.noregparm/s3vibes.o
usr/lib/oss/modules.noregparm/sblive.o
usr/lib/oss/modules.noregparm/sbxfi.o
usr/lib/oss/modules.noregparm/softoss.o
usr/lib/oss/modules.noregparm/solo.o
usr/lib/oss/modules.noregparm/trident.o
usr/lib/oss/modules.noregparm/via8233.o
usr/lib/oss/modules.noregparm/via97.o
usr/lib/oss/modules.noregparm/vortex.o
usr/lib/oss/modules.noregparm/ymf7xx.o
usr/lib/oss/modules.regparm/ali5455.o
usr/lib/oss/modules.regparm/allegro.o
usr/lib/oss/modules.regparm/als300.o
usr/lib/oss/modules.regparm/als4000.o
usr/lib/oss/modules.regparm/apci97.o
usr/lib/oss/modules.regparm/atiaudio.o
usr/lib/oss/modules.regparm/audigyls.o
usr/lib/oss/modules.regparm/audioloop.o
usr/lib/oss/modules.regparm/audiopci.o
usr/lib/oss/modules.regparm/cmi8788.o
usr/lib/oss/modules.regparm/cmpci.o
usr/lib/oss/modules.regparm/cs4280.o
usr/lib/oss/modules.regparm/cs4281.o
usr/lib/oss/modules.regparm/digi32.o
usr/lib/oss/modules.regparm/digi96.o
usr/lib/oss/modules.regparm/emu10k1x.o
usr/lib/oss/modules.regparm/envy24.o
usr/lib/oss/modules.regparm/envy24ht.o
usr/lib/oss/modules.regparm/fm801.o
usr/lib/oss/modules.regparm/geode.o
usr/lib/oss/modules.regparm/hdaudio.o
usr/lib/oss/modules.regparm/ich.o
usr/lib/oss/modules.regparm/imux.o
usr/lib/oss/modules.regparm/maestro.o
usr/lib/oss/modules.regparm/neomagic.o
usr/lib/oss/modules.regparm/ossusb.o
usr/lib/oss/modules.regparm/s3vibes.o
usr/lib/oss/modules.regparm/sblive.o
usr/lib/oss/modules.regparm/sbxfi.o
usr/lib/oss/modules.regparm/softoss.o
usr/lib/oss/modules.regparm/solo.o
usr/lib/oss/modules.regparm/trident.o
usr/lib/oss/modules.regparm/via8233.o
usr/lib/oss/modules.regparm/via97.o
usr/lib/oss/modules.regparm/vortex.o
usr/lib/oss/modules.regparm/ymf7xx.o
usr/lib/oss/objects
usr/lib/oss/objects.noregparm/osscore.o
usr/lib/oss/objects.regparm/osscore.o
usr/lib/oss/scripts/killprocs.sh
usr/lib/oss/scripts/remove_drv.sh
usr/lib/oss/scripts/restore_drv.sh
usr/lib/oss/scripts/setup-alsa.sh
usr/lib/oss/scripts/showprocs.sh
usr/lib/oss/sysfiles.list
usr/lib/oss/version.dat
usr/local/tce.installed/OSS
usr/local/tce.menu/OSS
usr/sbin/ossdetect
usr/sbin/ossdevlinks
usr/sbin/ossvermagic
usr/sbin/savemixer
usr/sbin/soundoff
usr/sbin/soundon
