# tcz/tc-ita.tcz.list
usr/local/bin/ab
usr/local/bin/add2file
usr/local/bin/appbrowser
usr/local/bin/appsaudit
usr/local/bin/cpanel
usr/local/bin/datetool
usr/local/bin/exittc
usr/local/bin/fdtool
usr/local/bin/filetool
usr/local/bin/filetool.sh
usr/local/bin/floppytool.sh
usr/local/bin/flpicsee
usr/local/bin/flrun
usr/local/bin/flwget
usr/local/bin/help
usr/local/bin/mktclocal
usr/local/bin/mnttool
usr/local/bin/mousetool
usr/local/bin/netcardconf
usr/local/bin/ondemand
usr/local/bin/select
usr/local/bin/stats
usr/local/bin/swapfile
usr/local/bin/tc-terminal-server
usr/local/bin/tce-audit
usr/local/bin/tce-load
usr/local/bin/tce-setdrive
usr/local/bin/tce-setup
usr/local/bin/tce-update
usr/local/bin/tcemirror
usr/local/bin/usbinstall
usr/local/bin/wallpaper
usr/local/bin/xsetup.sh
usr/local/etc/init.d/tc-functions
usr/local/etc/skel/.wmx/Appbrowser
usr/local/etc/skel/.wmx/Applications
usr/local/etc/skel/.wmx/Pannello_d_Controllo
usr/local/etc/skel/.wmx/Rinfresca_menu
usr/local/etc/skel/.wmx/Strumenti/Aggiungi_a_bootlocal
usr/local/etc/skel/.wmx/Strumenti/Aggiungi_a_filetool
usr/local/etc/skel/.wmx/Strumenti/Aggiungi_a_xfiletool
usr/local/etc/skel/.wmx/Strumenti/Editor d Testo
usr/local/etc/skel/.wmx/Strumenti/Imposta_Ora
usr/local/etc/skel/.wmx/Strumenti/Istantanea_Schermo
usr/local/etc/skel/.wmx/Strumenti/Processi_Attivi
usr/local/etc/skel/.wmx/Strumenti/TC_Terminal_Server
usr/local/etc/skel/.wmx/Strumenti/Uccidi_Finestra
usr/local/etc/skel/.wmx/Strumenti/VIsualizzatore_d'immagine
usr/local/etc/skel/.wmx/Strumenti/crea_tclocal
usr/local/etc/skel/.wmx/Terminali/Chiaro
usr/local/etc/skel/.wmx/Terminali/Scuro
usr/local/etc/skel/.wmx/Terminali/Super_Utente/Chiaro
usr/local/etc/skel/.wmx/Terminali/Super_Utente/Scuro
usr/local/etc/skel/.wmx/Terminali/Super_Utente/Trasparente
usr/local/etc/skel/.wmx/Terminali/Trasparente
usr/local/lib/locale/locale-archive
usr/local/tce.installed/tc-ita
